# deploy/systemd/stock-data-collector.service
[Unit]
Description=A-share objective news and market data collectors
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
User=__SERVICE_USER__
WorkingDirectory=__APP_DIR__
EnvironmentFile=/etc/stock-data.env
ExecStart=__APP_DIR__/venv/bin/python __APP_DIR__/scheduler.py --with-market --data-dir __APP_DIR__/data --collector-state __APP_DIR__/data/collector_health.json --market-output-dir __APP_DIR__/data_market --market-retention-days 0 --calendar-cache __APP_DIR__/data_market/trading_calendar.json --days 15
Restart=on-failure
RestartSec=20
NoNewPrivileges=true
PrivateTmp=true

[Install]
WantedBy=multi-user.target
